# services/krb5-server/supervisor.d/krb5.ini
[program:krb5kdc]
command = /usr/sbin/krb5kdc -n
autostart = true
autorestart = false
stdout_logfile = /var/log/krb5kdc.log
stdout_logfile_backups = 0
redirect_stderr = true

[program:kadmind]
command = /usr/sbin/kadmind -nofork
autostart = true
autorestart = false
stdout_logfile = /var/log/kadmind.log
stdout_logfile_backups = 0
redirect_stderr = true
